# freezer-display.service
[Unit]
Description=Freezer Inventory Display System
After=network.target graphical.target
Wants=graphical.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/freezer-inventory-system
ExecStart=/home/pi/freezer-inventory-system/start_pi_display.sh
Restart=always
RestartSec=10
Environment=DISPLAY=:0

[Install]
WantedBy=graphical.target
